# deploy/systemd/summariezer-digest.timer
[Unit]
Description=Run Summariezer Telegram digest every 3 days

[Timer]
OnActiveSec=3d
OnUnitActiveSec=3d
Persistent=true
Unit=summariezer-digest.service

[Install]
WantedBy=timers.target
