# deploy/systemd/summariezer-digest.service
[Unit]
Description=Summariezer Telegram digest job
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=root
Group=root
EnvironmentFile=/etc/summariezer/summariezer.env
WorkingDirectory=/opt/summariezer/app
ExecStart=/bin/bash /opt/summariezer/app/scripts/run_digest_window.sh
